# xinetd.d/nattcp
# default: off
# description: network performance testing server (developed from ttcp/nttcp/nuttcp)
service nattcp
{
	socket_type	= stream        
	wait		= no
	user		= nobody
	server		= /usr/bin/nattcp
	server_args	= -S
	log_on_failure	+= USERID
	disable		= yes
}

